# stremio.socket
[Unit]
Description=Socket for Stremio Docker App

[Socket]
ListenStream=443
Accept=yes

[Install]
WantedBy=sockets.target
